# lib/uksgx/external/dcap_source/tools/SGXPlatformRegistration/package/mpa_registration_tool.service
[Unit]
Description=Intel MPA Registration
Wants=network-online.target
After=network.target auditd.service network-online.target

[Service]
EnvironmentFile=/etc/environment
ExecStart=@mpa_folder@mpa_registration
[Install]
WantedBy=multi-user.target
