# lib/uksgx/external/dcap_source/tools/SGXPlatformRegistration/package/mpa_registration_tool.conf
description	"Intel(R) SGX Multi-Package Attestation Registration Service"

start on (local-filesystems and started networking)
stop on runlevel [!2345]
normal exit 0

env NAME="mpa_registration"
env MPA_PATH=@mpa_folder@

pre-start script
    test -x $MPA_PATH/$NAME || { stop; exit 0; }
end script
exec  $MPA_PATH/$NAME
