# lib/uksgx/external/dcap_source/QuoteGeneration/quote_wrapper/qgs/qgsd.conf
description	"Intel(R) TD Quoting Generation Service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
normal exit 0
umask 022

expect fork
console none
setuid "qgsd"

env NAME="qgs"
env QGS_PATH=@qgs_folder@
env DAEMON="$QGS_PATH/$NAME"


pre-start script
    test -x $QGS_PATH/$NAME || { stop; exit 0; }
end script


env LD_LIBRARY_PATH=$QGS_PATH
exec $QGS_PATH/$NAME
